# files/mitmproxy.service.j2
[Unit]
Description=mitmproxy

[Service]
Type=simple
ExecStart=/usr/local/bin/mitmproxy
Restart=always
RestartSec=10
User={{ mitmproxy.user }}

[Install]
WantedBy=multi-user.target
